# infra/production/etc/systemd/system/learning-app-backup-db.timer
[Unit]
Description=Learning App Backup DB Timer

[Timer]
OnCalendar=daily

[Install]
WantedBy=multi-user.target