# infra/production/etc/systemd/system/learning-app-dictionary-import.service
[Unit]
Description=Import dictionary data into CouchDB
After=network.target couchdb.service
Requires=couchdb.service

[Service]
Type=oneshot
User=webapp
WorkingDirectory=/var/lib/learning-app/dictionary
ExecStart=/opt/learning-app/dictionary/import.sh --db dictionary-db --reset

# Environment
EnvironmentFile=-/etc/learning-app/environment
EnvironmentFile=-/etc/environment.d/learning-app.conf

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes

[Install]
WantedBy=multi-user.target
