# infra/production/etc/systemd/system/learning-app-restart.path
[Unit]
Description=Learning App Restart Watch

[Path]
PathChanged=/opt/learning-app/learning-app.jar

[Install]
WantedBy=multi-user.target