# infra/production/DEBIAN/postinst
#!/usr/bin/env bash
set -euo pipefail

# Prompt admin for secrets and store securely
mkdir -p /etc/credstore.encrypted
if [ ! -f /etc/credstore.encrypted/borg-passphrase ]; then
    systemd-ask-password -n "Enter encryption password for BorgBackup" \
        | systemd-creds --name=borg-passphrase encrypt - /etc/credstore.encrypted/borg-passphrase
fi
if [ ! -f /etc/credstore.encrypted/openai_api_key ]; then
    systemd-ask-password -n "Enter API key for OpenAI" \
        | systemd-creds --name=openai_api_key encrypt - /etc/credstore.encrypted/openai_api_key
fi

# Create users
systemd-sysusers

# Create directories
systemd-tmpfiles --create

# Enable and start nginx config
nginx -t && systemctl reload nginx || systemctl start nginx

# Issue Let's Encrypt certificates interactively
if ! certbot certificates | grep -q sprecha.de; then
    certbot --nginx -d sprecha.de -d www.sprecha.de
fi

# Enable and start application services
systemctl daemon-reload
systemctl enable --now learning-app-run.service
systemctl enable --now learning-app-backup-db.timer
systemctl enable --now learning-app-restart.path
systemctl enable --now learning-app-certbot.timer
