# infra/production/DEBIAN/postinst
#!/usr/bin/env bash
set -euo pipefail

# Prompt admin for secrets and store securely
mkdir -p /etc/credstore.encrypted
if [ ! -f /etc/credstore.encrypted/borg-passphrase ]; then
    systemd-ask-password -n "Enter encryption password for BorgBackup" \
        | systemd-creds --name=borg-passphrase encrypt - /etc/credstore.encrypted/borg-passphrase
fi
if [ ! -f /etc/credstore.encrypted/openai_api_key ]; then
    systemd-ask-password -n "Enter API key for OpenAI" \
        | systemd-creds --name=openai_api_key encrypt - /etc/credstore.encrypted/openai_api_key
fi

# Create users
systemd-sysusers

# Create directories
systemd-tmpfiles --create

# Prepare nginx SSL directory and certbot deploy hook.
mkdir -p /etc/nginx/ssl/sprecha.de
chown root:root /etc/nginx/ssl/sprecha.de
chmod 700 /etc/nginx/ssl/sprecha.de

if [ -f /etc/letsencrypt/renewal-hooks/deploy/learning-app-nginx.sh ]; then
    chmod 700 /etc/letsencrypt/renewal-hooks/deploy/learning-app-nginx.sh
fi

# Enable and start nginx config
nginx -t && systemctl reload nginx || systemctl start nginx

# Issue Let's Encrypt certificates interactively
if ! certbot certificates | grep -q sprecha.de; then
    certbot --nginx -d sprecha.de -d www.sprecha.de
fi

# Ensure nginx uses root-only cert paths when certs exist.
if [ -f /etc/nginx/sites-available/learning-app.conf ]; then
    sed -i 's#/etc/letsencrypt/live/sprecha.de/fullchain.pem#/etc/nginx/ssl/sprecha.de/fullchain.pem#g' /etc/nginx/sites-available/learning-app.conf
    sed -i 's#/etc/letsencrypt/live/sprecha.de/privkey.pem#/etc/nginx/ssl/sprecha.de/privkey.pem#g' /etc/nginx/sites-available/learning-app.conf
fi

if [ -f /etc/letsencrypt/live/sprecha.de/fullchain.pem ]; then
    /etc/letsencrypt/renewal-hooks/deploy/learning-app-nginx.sh
fi

# Enable and start application services
systemctl daemon-reload
systemctl enable --now learning-app-run.service
systemctl enable --now learning-app-backup-db.timer
systemctl enable --now learning-app-restart.path
systemctl enable --now learning-app-certbot.timer
