# infra/etc/systemd/system/learning-app-restart.service
[Unit]
Description=Learning App Restart

[Service]
Type=oneshot
ExecStart=systemctl restart learning-app-run